# roles/backend/templates/backend.service.j2
[Unit]
Description=Backend Service

[Service]
User=Expense
Environment=DB_HOST={{ MYSQL_HOST }}
ExecStart=/usr/bin/node /app/index.js
SyslogIdentifier=backend

[Install]
WantedBy=multi-user.target